# deploy/cool-server/flick-update.timer
# ~/.config/systemd/user/flick-update.timer
[Unit]
Description=poll flick repos every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
